# low_battery.service
[Service]
ExecStart=/usr/lib/systemd/system/low_battery.d/low_battery.sh
User=video
Group=video
