# low_battery.timer
[Unit]
Description=Low battery monitor run on boot

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
